# flexget.cron
#!/bin/sh

exec s6-svc -o /services/flexget
