# services/flexget/run
#!/bin/sh

: ${FLEXGET_CONFIG:=data/config.yml}

exec 2>&1
exec flexget -c $FLEXGET_CONFIG execute
